# tags/asep/index.xml
<?xml version="1.0" encoding="utf-8" standalone="yes"?>
<rss version="2.0" xmlns:atom="http://www.w3.org/2005/Atom">
  <channel>
    <title>ASEP on Velociraptor - Digging deeper!</title>
    <link>https://docs.velociraptor.app/tags/asep/</link>
    <description>Recent content in ASEP on Velociraptor - Digging deeper!</description>
    <generator>Hugo</generator>
    <language>en-us</language>
    <lastBuildDate>Wed, 07 Jan 2026 23:39:40 -0500</lastBuildDate>
    <atom:link href="https://docs.velociraptor.app/tags/asep/index.xml" rel="self" type="application/rss+xml" />
    <item>
      <title>WMI Event Consumers: what are you missing?</title>
      <link>https://docs.velociraptor.app/blog/2022/2022-01-12-wmi-eventing/</link>
      <pubDate>Wed, 12 Jan 2022 00:00:00 +0000</pubDate>
      <guid>https://docs.velociraptor.app/blog/2022/2022-01-12-wmi-eventing/</guid>
      <description>&lt;p&gt;WMI Eventing is a fairly well known technique in DFIR, however some&#xA;tools may not provide the coverage you expect. This article covers&#xA;WMI eventing visibility and detection including custom namespaces.&lt;/p&gt;&#xA;&lt;p&gt;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&lt;figure&gt;&#xA;  &lt;div class=&#34;figure&#34;&gt;&#xA;    &lt;img src=&#34;https://docs.velociraptor.app/blog/2022/2022-01-12-wmi-eventing//blog/2022/2022-01-12-wmi-eventing/00SelectionBias.png&#34;&#xA;         class=&#34;inline&#34;&#xA;         alt=&#34;Selection bias in WWII: missing what is not collected.&#34;&#xA;         width=&#34;1000&#34; height=&#34;681&#34;&gt;&#xA;  &lt;/div&gt;&#xA;  &lt;figcaption&gt;&#xA;    &lt;a class=&#34;image-link&#34;&#xA;       target=&#34;_blank&#34; rel=&#34;noopener noreferrer&#34;&#xA;       href=&#34;00SelectionBias.png&#34;&gt;&#xA;      &lt;i class=&#34;fa fa-download&#34;&gt;&lt;/i&gt;&lt;/a&gt;&#xA;    Selection bias in WWII: missing what is not collected.&#xA;  &lt;/figcaption&gt;&#xA;&lt;/figure&gt;&#xA;&#xA;&#xA;&#xA;&#xA;&lt;/p&gt;&#xA;&lt;h2 id=&#34;background&#34;&gt;Background&lt;/h2&gt;&#xA;&lt;p&gt;There has been a fair bit of research and observations of WMI eventing&#xA;in field over the last years. In short, a WMI event consumer is a&#xA;method of subscribing to certain system events, then enabling an action&#xA;of some sort. Common adversary use cases may include persistence, privilege&#xA;escalation, or as a collection trigger. Represented as ATT&amp;amp;CK T1546.003&#xA;this technique has been observed in use from APT, through to trash-tic&#xA;worm and coin miner threats.&lt;/p&gt;&#xA;&lt;p&gt;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&lt;figure&gt;&#xA;  &lt;div class=&#34;figure&#34;&gt;&#xA;    &lt;img src=&#34;https://docs.velociraptor.app/blog/2022/2022-01-12-wmi-eventing//blog/2022/2022-01-12-wmi-eventing/01WMIOverview.png&#34;&#xA;         class=&#34;inline&#34;&#xA;         alt=&#34;WMI Eventing: 3 system classes&#34;&#xA;         width=&#34;1268&#34; height=&#34;614&#34;&gt;&#xA;  &lt;/div&gt;&#xA;  &lt;figcaption&gt;&#xA;    &lt;a class=&#34;image-link&#34;&#xA;       target=&#34;_blank&#34; rel=&#34;noopener noreferrer&#34;&#xA;       href=&#34;01WMIOverview.png&#34;&gt;&#xA;      &lt;i class=&#34;fa fa-download&#34;&gt;&lt;/i&gt;&lt;/a&gt;&#xA;    WMI Eventing: 3 system classes&#xA;  &lt;/figcaption&gt;&#xA;&lt;/figure&gt;&#xA;&#xA;&#xA;&#xA;&#xA;&lt;/p&gt;&#xA;&lt;p&gt;There are three system classes in every active event consumer:&lt;/p&gt;&#xA;&lt;ol&gt;&#xA;&lt;li&gt;__EventFilter is a WQL query that outlines the trigger event of&#xA;interest.&lt;/li&gt;&#xA;&lt;li&gt;__EventConsumer is an action to perform upon triggering an event.&lt;/li&gt;&#xA;&lt;li&gt;__FilterToConsumerBinding is the registration mechanism that binds&#xA;a filter to a consumer.&lt;/li&gt;&#xA;&lt;/ol&gt;&#xA;&lt;p&gt;Most detection will focus on collecting the WMI classes in root/subscription&#xA;and, in some tools root/default WMI namespaces.&lt;/p&gt;&#xA;&lt;p&gt;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&lt;figure&gt;&#xA;  &lt;div class=&#34;figure&#34;&gt;&#xA;    &lt;img src=&#34;https://docs.velociraptor.app/blog/2022/2022-01-12-wmi-eventing//blog/2022/2022-01-12-wmi-eventing/02Autoruns.png&#34;&#xA;         class=&#34;inline&#34;&#xA;         alt=&#34;Autoruns 14.07: detects root/default and root/subscription namespace WMI event consumers&#34;&#xA;         width=&#34;1056&#34; height=&#34;362&#34;&gt;&#xA;  &lt;/div&gt;&#xA;  &lt;figcaption&gt;&#xA;    &lt;a class=&#34;image-link&#34;&#xA;       target=&#34;_blank&#34; rel=&#34;noopener noreferrer&#34;&#xA;       href=&#34;02Autoruns.png&#34;&gt;&#xA;      &lt;i class=&#34;fa fa-download&#34;&gt;&lt;/i&gt;&lt;/a&gt;&#xA;    Autoruns 14.07: detects root/default and root/subscription namespace WMI event consumers&#xA;  &lt;/figcaption&gt;&#xA;&lt;/figure&gt;&#xA;&#xA;&#xA;&#xA;&#xA;&lt;/p&gt;&#xA;&lt;h4 id=&#34;custom-namespaces&#34;&gt;Custom Namespaces&lt;/h4&gt;&#xA;&lt;p&gt;At Blackhat 2018 Lee Christensen and Matt Graeber presented &amp;ldquo;Subverting&#xA;Sysmon: Application of a Formalized Security Product Evasion Methodology&amp;rdquo;.&#xA;This excellent talk focused on defense evasion methodology and highlighted&#xA;potential collection gaps in telemetry tools around WMI eventing. In this&#xA;case, the focus was on Sysmon behaviour of collection only in&#xA;root/subscription, interestingly, it also highlighted the possibility to&#xA;implement __EventConsumer classes in arbitrary namespaces.&lt;/p&gt;&#xA;&lt;p&gt;It is detection of WMI Event Consumers in arbitrary namespaces that I&amp;rsquo;m going&#xA;to focus. For anyone interested in testing I have written&#xA;&lt;a href=&#34;https://github.com/mgreen27/mgreen27.github.io/blob/master/static/other/WMIEventingNoisemaker/WmiEventingNoisemaker.ps1&#34; target=&#34;_blank&#34; &gt;a script to generate WMI event consumers&lt;/a&gt;.&#xA;This script wraps several powershell functions released during the Black&#xA;Hat talk to test creating working event consumers.&lt;/p&gt;&#xA;&lt;p&gt;First step was to create a custom namespace event consumer. In this&#xA;instance I selected the namespace name &lt;code&gt;totallylegit&lt;/code&gt; and attached an&#xA;ActiveScript event consumer.&lt;/p&gt;&#xA;&lt;p&gt;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&lt;figure&gt;&#xA;  &lt;div class=&#34;figure&#34;&gt;&#xA;    &lt;img src=&#34;https://docs.velociraptor.app/blog/2022/2022-01-12-wmi-eventing//blog/2022/2022-01-12-wmi-eventing/04WMIEventGeneration.png&#34;&#xA;         class=&#34;inline&#34;&#xA;         alt=&#34;WMIEventingNoismaker.ps1:Generate active script EventConsumer&#34;&#xA;         width=&#34;1048&#34; height=&#34;250&#34;&gt;&#xA;  &lt;/div&gt;&#xA;  &lt;figcaption&gt;&#xA;    &lt;a class=&#34;image-link&#34;&#xA;       target=&#34;_blank&#34; rel=&#34;noopener noreferrer&#34;&#xA;       href=&#34;04WMIEventGeneration.png&#34;&gt;&#xA;      &lt;i class=&#34;fa fa-download&#34;&gt;&lt;/i&gt;&lt;/a&gt;&#xA;    WMIEventingNoismaker.ps1:Generate active script EventConsumer&#xA;  &lt;/figcaption&gt;&#xA;&lt;/figure&gt;&#xA;&#xA;&#xA;&#xA;&#xA;&lt;/p&gt;&#xA;&lt;h2 id=&#34;collection&#34;&gt;Collection&lt;/h2&gt;&#xA;&lt;p&gt;Velociraptor has several valuable artifacts for hunting WMI Event&#xA;Consumers:&lt;/p&gt;&#xA;&lt;ul&gt;&#xA;&lt;li&gt;&lt;code&gt;Windows.Sysinternals.Autoruns&lt;/code&gt; - leverages a thirdparty deployment of&#xA;Sysinternals Autoruns and typically my go to ASEP collection artifact but&#xA;limited by visibility in root/default and root/subscription only.&lt;/li&gt;&#xA;&lt;li&gt;&lt;code&gt;Windows.Persistence.PermanentWMIEvents&lt;/code&gt; - recently upgraded to query&#xA;all ROOT namespaces.&lt;/li&gt;&#xA;&lt;/ul&gt;&#xA;&#xA;&lt;div class=&#34;mynotices tip&#34;&gt;&#xA;  &lt;div heading=&#34; Windows.Persistence.PermanentWMIEvents &#34;&gt;&lt;ul&gt;&#xA;&lt;li&gt;This artifact reports currently deployed permanent WMI Event Consumers.&lt;/li&gt;&#xA;&lt;li&gt;The artifact collects Binding information, then presents associated Filters and Consumers.&lt;/li&gt;&#xA;&lt;li&gt;Target a specific namespace, or tick &lt;code&gt;AllRootNamespaces&lt;/code&gt; to collect all&#xA;root namespace event consumers.&lt;/li&gt;&#xA;&lt;/ul&gt;&#xA;&lt;/div&gt;&#xA;&lt;/div&gt;&#xA;&#xA;&#xA;&lt;p&gt;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&lt;figure&gt;&#xA;  &lt;div class=&#34;figure&#34;&gt;&#xA;    &lt;img src=&#34;https://docs.velociraptor.app/blog/2022/2022-01-12-wmi-eventing//blog/2022/2022-01-12-wmi-eventing/05collection.png&#34;&#xA;         class=&#34;inline&#34;&#xA;         alt=&#34;Windows.Persistence.PermanentWMIEvents: configuration options&#34;&#xA;         width=&#34;1642&#34; height=&#34;330&#34;&gt;&#xA;  &lt;/div&gt;&#xA;  &lt;figcaption&gt;&#xA;    &lt;a class=&#34;image-link&#34;&#xA;       target=&#34;_blank&#34; rel=&#34;noopener noreferrer&#34;&#xA;       href=&#34;05collection.png&#34;&gt;&#xA;      &lt;i class=&#34;fa fa-download&#34;&gt;&lt;/i&gt;&lt;/a&gt;&#xA;    Windows.Persistence.PermanentWMIEvents: configuration options&#xA;  &lt;/figcaption&gt;&#xA;&lt;/figure&gt;&#xA;&#xA;&#xA;&#xA;&#xA;&lt;/p&gt;&#xA;&lt;p&gt;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&lt;figure&gt;&#xA;  &lt;div class=&#34;figure&#34;&gt;&#xA;    &lt;img src=&#34;https://docs.velociraptor.app/blog/2022/2022-01-12-wmi-eventing//blog/2022/2022-01-12-wmi-eventing/05collection_results.png&#34;&#xA;         class=&#34;inline&#34;&#xA;         alt=&#34;Windows.Persistence.PermanentWMIEvents: results&#34;&#xA;         width=&#34;1320&#34; height=&#34;492&#34;&gt;&#xA;  &lt;/div&gt;&#xA;  &lt;figcaption&gt;&#xA;    &lt;a class=&#34;image-link&#34;&#xA;       target=&#34;_blank&#34; rel=&#34;noopener noreferrer&#34;&#xA;       href=&#34;05collection_results.png&#34;&gt;&#xA;      &lt;i class=&#34;fa fa-download&#34;&gt;&lt;/i&gt;&lt;/a&gt;&#xA;    Windows.Persistence.PermanentWMIEvents: results&#xA;  &lt;/figcaption&gt;&#xA;&lt;/figure&gt;&#xA;&#xA;&#xA;&#xA;&#xA;&lt;/p&gt;&#xA;&lt;h4 id=&#34;telemetry&#34;&gt;Telemetry&lt;/h4&gt;&#xA;&lt;p&gt;Unfortunately prior to Windows 10 WMI logging was fairly limited. Sysmon and&#xA;other telemetry sources often rely on WMI eventing itself to collect WMI&#xA;eventing telemetry events. That means custom classes require namespace and&#xA;class existence prior to telemetry subscription. Sysmon as seen below also&#xA;does not have coverage for root/default namespace.&lt;/p&gt;&#xA;&lt;p&gt;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&lt;figure&gt;&#xA;  &lt;div class=&#34;figure&#34;&gt;&#xA;    &lt;img src=&#34;https://docs.velociraptor.app/blog/2022/2022-01-12-wmi-eventing//blog/2022/2022-01-12-wmi-eventing/03SysmonEid20.png&#34;&#xA;         class=&#34;inline&#34;&#xA;         alt=&#34;Sysmon collection: Event ID 20 mapping (__EventConsumer)&#34;&#xA;         width=&#34;1584&#34; height=&#34;722&#34;&gt;&#xA;  &lt;/div&gt;&#xA;  &lt;figcaption&gt;&#xA;    &lt;a class=&#34;image-link&#34;&#xA;       target=&#34;_blank&#34; rel=&#34;noopener noreferrer&#34;&#xA;       href=&#34;03SysmonEid20.png&#34;&gt;&#xA;      &lt;i class=&#34;fa fa-download&#34;&gt;&lt;/i&gt;&lt;/a&gt;&#xA;    Sysmon collection: Event ID 20 mapping (__EventConsumer)&#xA;  &lt;/figcaption&gt;&#xA;&lt;/figure&gt;&#xA;&#xA;&#xA;&#xA;&#xA;&lt;/p&gt;&#xA;&lt;p&gt;The good news is since Windows 10, WMI logging has improved significantly&#xA;and we can now query the event log: Microsoft-Windows-WMI-Activity or&#xA;subscribe the underlying ETW provider of the same name. In the VQL below&#xA;I filter the ETW event on event consumer creation or delete operations.&lt;/p&gt;&#xA;&lt;pre&gt;&lt;code class=&#34;language-vql&#34;&gt;SELECT&#xA;    System.TimeStamp AS EventTime,&#xA;    System.ID as EventId,&#xA;    strip(prefix=&#39;\\\\\.\\&#39;,string=EventData.NamespaceName) as NamespaceName,&#xA;    EventData.Operation as Operation,&#xA;    GetProcessInfo(TargetPid=int(int=EventData.ClientProcessId))[0] as Process&#xA;FROM watch_etw(guid=&amp;quot;{1418ef04-b0b4-4623-bf7e-d74ab47bbdaa}&amp;quot;)&#xA;WHERE EventId = 11&#xA;    AND Operation =~ &#39;WbemServices::(PutInstance|DeleteInstance|PutClass|DeleteClass)&#39;&#xA;    AND Operation =~ &#39;EventConsumer|EventFilter|FilterToConsumerBinding&#39;&#xA;&lt;/code&gt;&lt;/pre&gt;&#xA;&lt;p&gt;I have included a completed artifact in the artifact exchange:&#xA;&lt;a href=&#34;https://docs.velociraptor.app/exchange/artifacts/pages/wmieventing/&#34;&gt;Windows.ETW.WMIEventing&lt;/a&gt;.&#xA;That artifact includes process enrichment, targeting both creation and deletion of EventConsumers.&lt;/p&gt;&#xA;&lt;p&gt;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&lt;figure&gt;&#xA;  &lt;div class=&#34;figure&#34;&gt;&#xA;    &lt;img src=&#34;https://docs.velociraptor.app/blog/2022/2022-01-12-wmi-eventing//blog/2022/2022-01-12-wmi-eventing/06ETW.png&#34;&#xA;         class=&#34;inline&#34;&#xA;         alt=&#34;Custom namespace provider registration and process enrichment&#34;&#xA;         width=&#34;1129&#34; height=&#34;480&#34;&gt;&#xA;  &lt;/div&gt;&#xA;  &lt;figcaption&gt;&#xA;    &lt;a class=&#34;image-link&#34;&#xA;       target=&#34;_blank&#34; rel=&#34;noopener noreferrer&#34;&#xA;       href=&#34;06ETW.png&#34;&gt;&#xA;      &lt;i class=&#34;fa fa-download&#34;&gt;&lt;/i&gt;&lt;/a&gt;&#xA;    Custom namespace provider registration and process enrichment&#xA;  &lt;/figcaption&gt;&#xA;&lt;/figure&gt;&#xA;&#xA;&#xA;&#xA;&#xA;&lt;/p&gt;&#xA;&lt;p&gt;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&lt;figure&gt;&#xA;  &lt;div class=&#34;figure&#34;&gt;&#xA;    &lt;img src=&#34;https://docs.velociraptor.app/blog/2022/2022-01-12-wmi-eventing//blog/2022/2022-01-12-wmi-eventing/06ETWb.png&#34;&#xA;         class=&#34;inline&#34;&#xA;         alt=&#34;Windows.ETW.WMIEventing: all operations event consumer creation and removal&#34;&#xA;         width=&#34;863&#34; height=&#34;301&#34;&gt;&#xA;  &lt;/div&gt;&#xA;  &lt;figcaption&gt;&#xA;    &lt;a class=&#34;image-link&#34;&#xA;       target=&#34;_blank&#34; rel=&#34;noopener noreferrer&#34;&#xA;       href=&#34;06ETWb.png&#34;&gt;&#xA;      &lt;i class=&#34;fa fa-download&#34;&gt;&lt;/i&gt;&lt;/a&gt;&#xA;    Windows.ETW.WMIEventing: all operations event consumer creation and removal&#xA;  &lt;/figcaption&gt;&#xA;&lt;/figure&gt;&#xA;&#xA;&#xA;&#xA;&#xA;&lt;/p&gt;&#xA;&lt;h4 id=&#34;event-log&#34;&gt;Event Log&lt;/h4&gt;&#xA;&lt;p&gt;Similar filters can be used with &lt;code&gt;Windows.EventLogs.EvtxHunter&lt;/code&gt; for&#xA;detection. Its worthy to note, event logs hold less verbose logging for&#xA;the registration than ETW but this use case is helpful when coming late&#xA;to the party during an investigation.&lt;/p&gt;&#xA;&lt;p&gt;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&lt;figure&gt;&#xA;  &lt;div class=&#34;figure&#34;&gt;&#xA;    &lt;img src=&#34;https://docs.velociraptor.app/blog/2022/2022-01-12-wmi-eventing//blog/2022/2022-01-12-wmi-eventing/07EvtxHunter.png&#34;&#xA;         class=&#34;inline&#34;&#xA;         alt=&#34;Windows.EventLogs.EvtxHunter: hunt for event consumer string&#34;&#xA;         width=&#34;472&#34; height=&#34;363&#34;&gt;&#xA;  &lt;/div&gt;&#xA;  &lt;figcaption&gt;&#xA;    &lt;a class=&#34;image-link&#34;&#xA;       target=&#34;_blank&#34; rel=&#34;noopener noreferrer&#34;&#xA;       href=&#34;07EvtxHunter.png&#34;&gt;&#xA;      &lt;i class=&#34;fa fa-download&#34;&gt;&lt;/i&gt;&lt;/a&gt;&#xA;    Windows.EventLogs.EvtxHunter: hunt for event consumer string&#xA;  &lt;/figcaption&gt;&#xA;&lt;/figure&gt;&#xA;&#xA;&#xA;&#xA;&#xA;&lt;/p&gt;&#xA;&lt;p&gt;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&#xA;&lt;figure&gt;&#xA;  &lt;div class=&#34;figure&#34;&gt;&#xA;    &lt;img src=&#34;https://docs.velociraptor.app/blog/2022/2022-01-12-wmi-eventing//blog/2022/2022-01-12-wmi-eventing/07EvtxHunterb.png&#34;&#xA;         class=&#34;inline&#34;&#xA;         alt=&#34;Windows.EventLogs.EvtxHunter: detect event consumer class creation&#34;&#xA;         width=&#34;667&#34; height=&#34;673&#34;&gt;&#xA;  &lt;/div&gt;&#xA;  &lt;figcaption&gt;&#xA;    &lt;a class=&#34;image-link&#34;&#xA;       target=&#34;_blank&#34; rel=&#34;noopener noreferrer&#34;&#xA;       href=&#34;07EvtxHunterb.png&#34;&gt;&#xA;      &lt;i class=&#34;fa fa-download&#34;&gt;&lt;/i&gt;&lt;/a&gt;&#xA;    Windows.EventLogs.EvtxHunter: detect event consumer class creation&#xA;  &lt;/figcaption&gt;&#xA;&lt;/figure&gt;&#xA;&#xA;&#xA;&#xA;&#xA;&lt;/p&gt;&#xA;&lt;h1 id=&#34;conclusions&#34;&gt;Conclusions&lt;/h1&gt;&#xA;&lt;p&gt;During this post, we have shown three techniques for detecting WMI event consumers&#xA;that are worth considering. We can collect these data-points over an entire&#xA;network in minutes using Velociraptor&amp;rsquo;s &amp;ldquo;hunt&amp;rdquo; capability. Similarly&#xA;Velociraptor notebook workflow assists excluding known good entries quickly as part of analysis.&lt;/p&gt;&#xA;&lt;p&gt;The Velociraptor platform aims to provide visibility and access&#xA;to endpoint data. If you would like to try Velociraptor it is available on GitHub under an open source license.&#xA;As always, please file issues on the bug tracker or ask questions on our&#xA;mailing list &lt;a href=&#34;mailto:velociraptor-discuss@googlegroups.com&#34; target=&#34;_blank&#34; &gt;velociraptor-discuss@googlegroups.com&lt;/a&gt;. You can also chat with&#xA;us directly on discord at &lt;a href=&#34;https://www.velocidex.com/discord&#34; target=&#34;_blank&#34; &gt;https://www.velocidex.com/discord&lt;/a&gt;&lt;/p&gt;&#xA;&lt;h2 id=&#34;references&#34;&gt;References&lt;/h2&gt;&#xA;&lt;ol&gt;&#xA;&lt;li&gt;&lt;a href=&#34;https://docs.microsoft.com/en-us/windows/win32/wmisdk/about-wmi&#34; target=&#34;_blank&#34; &gt;Microsoft documentation, About WMI&lt;/a&gt;&lt;/li&gt;&#xA;&lt;li&gt;&lt;a href=&#34;https://attack.mitre.org/techniques/T1546/003/&#34; target=&#34;_blank&#34; &gt;MITRE ATT&amp;amp;CK T1546.003, Event Triggered Execution: Windows Management Instrumentation Event Subscription&lt;/a&gt;&lt;/li&gt;&#xA;&lt;li&gt;&lt;a href=&#34;https://www.youtube.com/watch?v=R5IEyoFpZq0&#34; target=&#34;_blank&#34; &gt;Christensen.L and Graeber.M, Blackhat 2018 - Subverting Sysmon: Application of a Formalized Security Product Evasion Methodology&lt;/a&gt;&lt;/li&gt;&#xA;&lt;li&gt;&lt;a href=&#34;https://github.com/jsecurity101/Windows-API-To-Sysmon-Events/&#34; target=&#34;_blank&#34; &gt;JSecurity101, Windows APIs To Sysmon-Events&lt;/a&gt;&lt;/li&gt;&#xA;&lt;/ol&gt;&#xA;</description>
    </item>
  </channel>
</rss>
